# applications/mongot/s6-config/s6-rc.d/mongod/finish
#!/bin/sh
echo "[mongod] process exited with code $1"
