# applications/mongot/s6-config/s6-rc.d/mongod/run
#!/bin/sh
exec mongod --config /etc/mongod.conf 2>&1 \
  | awk '{ print "[mongod] " $0; fflush() }'
