# applications/mongot/s6-config/s6-rc.d/mongot/finish
#!/bin/sh
echo "[mongot] process exited with code $1"
